# service-dashboard/sysdaemon/sysdaemon.conf
description "sysdaemon startscript"

start on runlevel [2345]
stop on runlevel [!2345]

exec python3 /opt/optimistic-failover/service-dashboard/sysdaemon/sysdaemon.py > /var/log/sysdaemon.log
